# spark-cluster/docker/resources/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-cluster]
command=/opt/init.sh -DFOREGROUND

